# deploy/gcp/systemd/xriq-indexer.service
[Unit]
Description=XRIQ staging-devnet indexer -> Cloud SQL (one-shot; run by the timer)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/xriq-indexer-run.sh
