# deploy/gcp/systemd/xriq-testnet-node.service
[Unit]
Description=XRIQ testnet node (serve-readonly --network testnet, TEST-ONLY, no value)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=-/usr/bin/docker rm -f xriq-testnet-node
ExecStart=/usr/local/bin/xriq-testnet-node-run.sh
ExecStop=-/usr/bin/docker stop xriq-testnet-node
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
